# display-management/systemd/wlr-display-off-nightly@.timer
[Unit]
Description=Turn OFF output %i nightly

[Timer]
OnCalendar=22:00
Persistent=true

[Install]
WantedBy=timers.target
